# templates/etc/systemd/system/docker.service.d/options.conf.j2
# {{ ansible_managed }}

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd {{ boss__docker__daemon_flags | join(" ") }}
